# systemd/ceph-api.service
[Unit]
Description=Ceph Management REST API
Documentation=https://github.com/your-org/cephx-api
After=network-online.target ceph.target
Wants=network-online.target

[Service]
Type=notify
User=zenjabba
Group=zenjabba
WorkingDirectory=/opt/ceph-api

# Environment
Environment="PATH=/opt/ceph-api/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=-/etc/default/ceph-api

# Main service execution
ExecStart=/opt/ceph-api/venv/bin/uvicorn app.main:app \
    --host 10.10.3.4 \
    --port 8080 \
    --workers 4 \
    --log-level info \
    --access-log \
    --proxy-headers \
    --forwarded-allow-ips '*'

# Restart policy
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ceph-api/data /var/log/ceph-api
ReadOnlyPaths=/etc/ceph

# Capability restrictions
CapabilityBoundingSet=
AmbientCapabilities=

# System call filtering
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
SystemCallErrorNumber=EPERM

# Namespace isolation
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
RestrictSUIDSGID=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=512
MemoryMax=1G
TasksMax=512
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ceph-api

[Install]
WantedBy=multi-user.target
